# mvs_openvino_demo.2/mvs-detector.service
[Unit]
Description=迈德威视 MVS OpenVINO 实时目标检测
After=graphical-session.target display-manager.service
Wants=graphical-session.target

[Service]
Type=simple

# 项目路径（根据 CMakeLists.txt: project(mvs_openvino_demo)）
WorkingDirectory=%h/mvs_openvino_demo.1
ExecStart=%h/mvs_openvino_demo.1/start.sh

# 环境变量
Environment="DISPLAY=:0"

# ★ 退出后自动重启（通信流程结束后程序会退出，需要自动重启进入下一轮）
Restart=always
RestartSec=10

# 启动超时（模型加载 + 预热需要时间）
TimeoutStartSec=120

# 优雅退出
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mvs-detector

# 降低 OOM 被杀概率
OOMScoreAdjust=-200

[Install]
WantedBy=default.target
